# crontab
# Nightly backup to Synology NAS
# Runs at 2:00 AM every day
0 2 * * * $HOME/bin/backup
