# systemd/system/user@.service.d/20-oomd-app-slice.conf
# Application pressure is monitored within the user's app.slice.
[Service]
ManagedOOMMemoryPressure=auto
ManagedOOMMemoryPressureLimit=0%
